# config/supervisor/squash-web.production.conf
[program:squash-web]
command=/home/squash-web/.rvm/bin/rvm jruby-1.7.13@squash-web do bundle exec rails s -s trinidad -p 3002
directory=/home/squash-web/current/
user=squash-web
autostart=true
autorestart=true
startsecs=2
startretries=100
stdout_logfile=/home/squash-web/current/log/production.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=50
stderr_logfile=/home/squash-web/current/log/production.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=50
stopsignal=TERM
environment=RAILS_ENV=production